# supervisord.conf
[supervisord]
nodaemon = true
user = root

[program:homeassistant]
command = python -m homeassistant --verbose --config /config

[program:telldusd]
command = telldusd --nodaemon
