# supervisord.conf
[supervisord]
nodaemon=true

[program:homeassistant]
command=python -m homeassistant --config /config

[program:telldusd]
command=telldusd --nodaemon
